# roles/api-taxi/templates/systemd/api-taxi.service.j2
#
# This file is managed by ansible. Do not edit manually.
#

[Unit]

Description=api-taxi
BindsTo=docker.service
After=docker.service

[Service]

Restart=on-failure
RestartSec=5

# Pulling the image might be long in case of slow network
TimeoutStartSec=20min

# Pull latest version on start
ExecStartPre=/usr/bin/docker pull {{ api_taxi_image }}:{{ api_taxi_version }}

# Remove existing containers, if exist
ExecStartPre=-/usr/bin/docker kill api-taxi
ExecStartPre=-/usr/bin/docker rm api-taxi

# Run container
ExecStart=/usr/bin/docker run --name api-taxi \
    -v /etc/openmaraude/api-taxi.ini:/etc/uwsgi-emperor/vassals/api-taxi.ini:ro \
    -v /var/run/openmaraude:/var/run/openmaraude \
    openmaraude/apitaxi:latest

# Stop container
ExecStop=/usr/bin/docker stop api-taxi

# Send SIGHUP to uwsgi-emperor
ExecReload=/usr/bin/docker kill -s SIGHUP api-taxi

[Install]
WantedBy=multi-user.target
